# src/main/doc/crontab
crontab -e

*/10 * * * * sh /root/updateIp.sh